# Durex.service
[Unit]
Description=NULL

[Service]
Type=forking
ExecStart=/sbin/Durex

[Install]
WantedBy=multi-user.target

